# pkg_scripts/spinnaker-example-go.conf
description "spinnaker-example-go"

expect fork

start on runlevel [2345]
stop on runlevel [!2345]

exec /opt/spinnaker-example-go 2>&1 > /opt/spinnaker-example-go.log &
